# raspi_zero_2_w_code/streamer.service
[Unit]
Description=enel Pi RTSP streamer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/streamer.sh
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Pi Zero 2 W: allow picamera/libcamera + RTSP networking
User=root
Group=root

[Install]
WantedBy=multi-user.target
