# raspi_zero_2_w_code/wifi-j-autojoin.timer
[Unit]
Description=enel: periodic hotspot/J Wi-Fi scan

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
Unit=wifi-j-autojoin.service

[Install]
WantedBy=timers.target
