# raspi_zero_2_w_code/wifi-j-autojoin.service
[Unit]
Description=enel: auto-join operator hotspot or open J Wi-Fi
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-j-autojoin.sh

[Install]
WantedBy=multi-user.target
